# multiple_project/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn-{{ project_name }}.sock

[Install]
WantedBy=sockets.target